# ansible/roles/loki/templates/loki.service.j2
[Unit]
Description=Loki logging service

[Service]
ExecStart=/usr/local/bin/loki -config.file=/etc/loki/config.yml

[Install]
WantedBy=multi-user.target